# gocollect-client/rc/upstart.conf
# vim: set syn=upstart:
description "GoCollect, the rarely changing sysinfo collector"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec gocollect
