# rc/upstart.conf
# vim: set syn=upstart:
description "GoCollect collects system info and pushes it to a central server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec gocollect
